# crontab
*/1 * * * * python /cloudflare-api-request.py >> /logfile
